# src/miners/rss/tracker-miner-rss.service.in
[Unit]
Description=Tracker RSS/ATOM feed data miner

[Service]
Type=dbus
BusName=org.freedesktop.Tracker1.Miner.RSS
ExecStart=@libexecdir@/tracker-miner-rss
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
RestartPreventExitStatus=SIGKILL
