# src/miners/user-guides/tracker-miner-user-guides.service.in
[Unit]
Description=Tracker user guides data miner

[Service]
Type=dbus
BusName=org.freedesktop.Tracker1.Miner.UserGuides
ExecStart=@libexecdir@/tracker-miner-user-guides
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
RestartPreventExitStatus=SIGKILL
