# src/tracker-store/tracker-store.service.in
[Unit]
Description=Tracker metadata database store and lookup manager

[Service]
Type=dbus
BusName=org.freedesktop.Tracker1
ExecStart=@libexecdir@/tracker-store
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
RestartPreventExitStatus=SIGKILL
