# src/tracker-extract/tracker-extract.service.in
[Unit]
Description=Tracker metadata extractor

[Service]
Type=dbus
BusName=org.freedesktop.Tracker1.Miner.Extract
ExecStart=@libexecdir@/tracker-extract
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
RestartPreventExitStatus=SIGKILL
