# src/portal/tracker-xdg-portal.service.in
[Unit]
Description=Tracker portal service

[Service]
Type=dbus
BusName=org.freedesktop.portal.Tracker
ExecStart=@libexecdir@/tracker-xdg-portal-3
