# src/portal/tracker-xdg-portal.service.in
[Unit]
Description=TinySPARQL portal service

[Service]
Type=dbus
BusName=org.freedesktop.portal.Tracker
ExecStart=@libexecdir@/tinysparql-xdg-portal-3
