# src/tracker-writeback/tracker-writeback.service.in
[Unit]
Description=Tracker data writeback service

[Service]
Type=dbus
BusName=org.freedesktop.Tracker1.Writeback
ExecStart=@libexecdir@/tracker-writeback
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
RestartPreventExitStatus=SIGKILL
