# src/miners/apps/tracker-miner-apps.service.in
[Unit]
Description=Tracker applications data miner

[Service]
Type=dbus
BusName=org.freedesktop.Tracker1.Miner.Applications
ExecStart=@libexecdir@/tracker-miner-apps
Restart=on-failure
# Don't restart after tracker daemon -k (aka tracker-control -k)
RestartPreventExitStatus=SIGKILL
